# crontab
SHELL="/bin/bash"
HOME="/home/dylan"
PATH="/bin:/usr/bin:$HOME/bin/bin"

1 0 * * * /bin/bash /home/dylan/.files/bin/update_linux_repo.sh &> /home/dylan/linux-git.log;

